# lecture-11-13/intro/target/a.service
[Unit]
Description= Test Service A


[Service]
ExecStart=sleep 2
Type=oneshot
RemainAfterExit=true


[Install]
WantedBy=test.target
