# lecture-11-13/intro/target/b.service
[Unit]
Description= Test Service B
After=a.service
#Wants=a.service
Requires=a.service

[Service]
ExecStart=sleep 2
Type=oneshot
RemainAfterExit=true

[Install]
WantedBy=test.target
