# lecture-11-15/intro2/target/a.service
[Unit]
Description= Test Service A
Conflicts=stop.target

[Service]
ExecStart=sleep 2
Type=oneshot
RemainAfterExit=true


[Install]
WantedBy=test.target
