# lecture-11-15/intro2/target/b.service
[Unit]
Description= Test Service B
After=a.service
Conflicts=stop.target
#Wants=a.service
Requires=a.service

[Service]
ExecStart=sleep 2
Type=oneshot
RemainAfterExit=true

[Install]
WantedBy=test.target
